# cookbooks/wikijs/templates/default/wikijs.ini.erb
[program:wikijs]
command = node server
autostart=true
autorestart=true
priority=5
directory=<%= node['wikijs']['installDir'] %>
